# src/server/openrc/nymphcast
#!/sbin/openrc-run

supervisor=supervise-daemon
command="/opt/nymphcast/nymphcast_server"
command_args="-c /opt/nymphcast/nymphcast_config.ini"
