# src/server/openrc/nymphcast
#!/sbin/openrc-run

supervisor=supervise-daemon
command="@BIN@"
command_args="-c @CONFIG@ -a @APPS@ -w @WALLPAPERS@"
