# src/server/systemd/nymphcast.service
[Unit]
Description=NymphCast audio/video casting server

[Service]
ExecStart=/opt/nymphcast/nymphcast_server -c nymphcast_video_config.ini
User=pi
WorkingDirectory=/opt/nymphcast
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target